# projects/Rockchip/devices/RG351P/filesystem/usr/lib/systemd/system/aftersuspend.service
[Unit]
Description=After Suspend
After=suspend.target

[Service]
Type=oneshot
ExecStart=-amixer set "Playback Path" "SPK_HP"

[Install]
WantedBy=suspend.target
