# packages/sysutils/busybox/profile.d/098-busybox
# SPDX-License-Identifier: GPL-2.0-or-later
# Copyright (C) 2009-2014 Stephan Raue (stephan@openelec.tv)
# Copyright (C) 2021-present Fewtarius

. /etc/os-release

export HOME="/storage"
export PATH="/usr/bin:/usr/sbin"
export HOSTNAME=`cat /etc/hostname`
export PS1='\e[1;30m\h:\e[1;30m\w \e[0m\$ \e[0m'

# Redirect requests for sudo to nowhere.
alias sudo=''

case "$TERM" in
  linux|nxterm|screen)
    ;;
  "")
    export TERM="linux"
    ;;
  *)
    export TERM="xterm"
    ;;
esac
