# packages/ui/emulationstation/system.d/emustation.service
[Unit]
Description=EmulationStation emulator frontend
ConditionPathExists=/var/lock/start.games
After=jelos-automount.service

[Service]
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=HOME=/storage
Environment=SDL_AUDIODRIVER=pulseaudio
ExecStart=/usr/bin/start_es.sh
TimeoutStopSec=3
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=jelos.target
