# packages/ui/emulationstation/system.d/emustation.service
[Unit]
Description=EmulationStation emulator frontend
ConditionPathExists=/var/lock/start.games

[Service]
Environment=HOME=/storage
Environment=SDL_AUDIODRIVER=pulseaudio
ExecStartPre=/usr/bin/es_settings
ExecStart=/usr/bin/emulationstation --log-path /var/log --no-splash
KillMode=process
TimeoutStopSec=3
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=jelos.target
