# packages/ui/emulationstation/system.d/emustation.service
[Unit]
Description=EmulationStation emulator frontend
ConditionPathExists=/var/lock/start.games

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStartPre=/usr/bin/emustation-config
ExecStart=/usr/bin/emulationstation --log-path /var/log
KillMode=process
TimeoutStopSec=3
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=jelos.target
