# packages/ui/emulationstation/system.d/essway.service
[Unit]
Description=EmulationStation-with-Sway
# This path never exist on rk3566
#ConditionPathExists=/var/lock/start.games
Requires=sway.service
After=sway.service


[Service]
Environment=HOME=/storage
Type=simple
ExecStart=/usr/bin/start_es.sh
# This is  hack, ES doesn't immediately detect the builtin gamepad
# unless an input device is connected while ES is running
ExecStartPost=/usr/bin/timeout --preserve-status 1 /usr/bin/gptokeyb
WorkingDirectory=/storage
Restart=always
RestartSec=2


[Install]
WantedBy=jelos.target
