# packages/apps/jelos-gamepad/system.d/jelos_gamepad.service
[Unit]
Description=Ragnarok Input Daemon

[Service]
Type=simple
ExecStartPre=/usr/bin/udevadm settle
ExecStart=/usr/bin/jelos_gamepad

[Install]
WantedBy=multi-user.target
