# packages/network/zerotier-one/system.d/zerotier-one.service
[Unit]
Description=zerotier-one
Wants=network-pre.target
After=network-pre.target

[Service]
ExecStart=/usr/sbin/zerotier-one
Restart=always
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
