# packages/jelos/system.d/bluetooth-agent.service
[Unit]
Description=Start bluetooth agent
After=bluetooth.service
PartOf=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/batocera-bluetooth-agent

[Install]
WantedBy=bluetooth.service
