# packages/jelos/system.d/jelos-automount.service
[Unit]
Description=JELOS user autostart script
Before=autostart.service
After=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/automount
ExecStartPost=systemd-tmpfiles --create /usr/lib/tmpfiles.d/jelos-dirs.conf
RemainAfterExit=yes

[Install]
WantedBy=jelos.target
