# packages/jelos/system.d/save-sysconfig.service
[Unit]
Description=Back up the system config before shutting down
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/bin/chksysconfig backup
TimeoutStartSec=0

[Install]
WantedBy=shutdown.target
