# packages/network/simple-http-server/system.d/simple-http-server.service
[Unit]
Description=simple-http-server
Wants=network-pre.target
After=network-pre.target

[Service]
ExecStart=/usr/bin/bash -c "source /etc/profile.d/001-functions && /usr/bin/simple-http-server -i -p 80 -a root:$(get_setting root.password) -s -u /storage"
Restart=always
KillMode=process
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
