# packages/emulators/tools/virtualgamepad/system.d/virtualgamepad.service
[Unit]
Description=Virtual Gamepad Driver

[Service]
Type=simple
ExecStart=/usr/bin/virtualgamepad

[Install]
WantedBy=multi-user.target
