# packages/sysutils/system-utils/system.d/headphones.service
[Unit]
Description=Headphone Sense
Before=jelos.target

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/headphone_sense
Restart=on-failure

[Install]
WantedBy=multi-user.target
