# packages/sysutils/system-utils/system.d/video.service
[Unit]
Description=Video Sense
Before=jelos.target

[Service]
Environment=HOME=/storage
EnvironmentFile=/etc/profile
ExecStart=/usr/bin/video_sense
Restart=on-failure

[Install]
WantedBy=multi-user.target
