# packages/virtual/emulators/system.d/storage-roms.mount
[Unit]
Description=Overlays directory
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/storage/roms
Type=overlay
Options=lowerdir=/storage/games-external,upperdir=/storage/games-internal,workdir=/storage/.tmp/games-workdir

[Install]
WantedBy=jelos.target
