# packages/sysutils/util-linux/system.d/fstrim.service
[Unit]
Description=Discard unused blocks on all filesystems

[Service]
Type=oneshot
ExecStart=-/usr/sbin/fstrim --all --verbose --quiet-unsupported
ExecStart=-/usr/bin/touch /storage/.cache/fstrim.run
