# packages/emulators/tools/virtualcontroller/system.d/virtualcontroller.service
[Unit]
Description=Virtual Controller Driver
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
Nice=-20
OOMScoreAdjust=-1000
Restart=always
ExecStart=/usr/bin/virtual_controller

[Install]
WantedBy=multi-user.target
