# packages/audio/pipewire/system.d/pipewire-pulse.socket
[Unit]
Description=PipeWire PulseAudio
Conflicts=pulseaudio.socket

[Socket]
Priority=6
ListenStream=%t/pulse/native

[Install]
WantedBy=sockets.target
