# packages/audio/bluez-alsa/system.d/bluealsa-defaults.service
[Unit]
Description=Bluealsa defaults
After=local-fs.target

ConditionPathExists=!/storage/.cache/services/bluealsa.conf
ConditionPathExists=!/storage/.cache/services/bluealsa.disabled

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cp /usr/share/services/bluealsa.conf /storage/.cache/services/'
RemainAfterExit=yes
