# packages/sysutils/systemd/system.d/userconfig.service
[Unit]
Description=Setup User config dir
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
Before=jelos-automount.service autostart.service

[Service]
Type=oneshot
ExecStart=/usr/bin/userconfig-setup
RemainAfterExit=yes
StandardOutput=journal+console

[Install]
WantedBy=sysinit.target
