# packages/jelos/profile.d/98-jslisten
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
# Copyright (C) 2021-present Fewtarius (https://github.com/fewtarius)

. /etc/os-release

export PATH="$PATH:/usr/local/bin:/usr/bin:/storage/bin"

jslisten() {
        JSLISTENCONF="/storage/.config/system/configs/jslisten.cfg"
        if [ "$1" == "set" ]
        then
                systemctl stop jslisten
		rm /storage/.config/system/configs/jslisten.cfg

                if [ "$2" == "mpv" ]
                then

                        cat <<EOF >${JSLISTENCONF}
[Pause]
program="/usr/bin/manage_mpv.sh pause"
button1=${DEVICE_BTN_START}

[Exit]
program="/usr/bin/manage_mpv.sh quit"
button1=${DEVICE_BTN_SELECT}
button2=${DEVICE_BTN_START}

[Skip5s]
program="/usr/bin/manage_mpv.sh skip5s"
button1=${DEVICE_BTN_TR}

[Back5s]
program="/usr/bin/manage_mpv.sh back5s"
button1=${DEVICE_BTN_TL}

[Skip60s]
program="/usr/bin/manage_mpv.sh skip60s"
button1=${DEVICE_BTN_TR2}

[Back60s]
program="/usr/bin/manage_mpv.sh back60s"
button1=${DEVICE_BTN_TL2}

EOF
		fi

                cat <<EOF >>${JSLISTENCONF}
[KillAll]
program="/usr/bin/killall ${2}"
button1=${DEVICE_BTN_TL}
button2=${DEVICE_BTN_SELECT}
button3=${DEVICE_BTN_START}

[BrightnessUp]
program="/usr/bin/brightness up"
button1=${DEVICE_BTN_TL}
button2=${DEVICE_BTN_DPAD_UP}

[BrightnessDown]
program="/usr/bin/brightness down"
button1=${DEVICE_BTN_TL}
button2=${DEVICE_BTN_DPAD_DOWN}

EOF
                systemctl start jslisten
        elif [ "$1" == "stop" ]
        then
                systemctl stop jslisten
        elif [ "$1" == "start" ]
        then
                systemctl start jslisten
        fi
}
