# packages/jelos/profile.d/99-freqfunctions
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
# Copyright (C) 2021-present Fewtarius (https://github.com/fewtarius)

. /etc/os-release

set_gov() {
  for POLICY in ${CPU_FREQ[@]}
  do
    echo $1 >${POLICY}/scaling_governor 2>/dev/null
  done

  for HANDLE in ${GPU_FREQ} ${DMC_FREQ}
  do
    echo $1 >${HANDLE}/governor 2>/dev/null || echo dmc_$1 >${HANDLE}/governor 2>/dev/null
  done
}

performance() {
  set_gov performance
}

ondemand() {
  set_gov ondemand
}

interactive() {
  set_gov interactive
}

powersave() {
  set_gov powersave
}

