# packages/jelos/profile.d/99-freqfunctions
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
# Copyright (C) 2021-present Fewtarius (https://github.com/fewtarius)

. /etc/os-release


set_cpu_gov() {
  for POLICY in ${CPU_FREQ[@]}
  do
    echo $1 >${POLICY}/scaling_governor 2>/dev/null
  done
}

set_dmc_gov() {
  echo $1 >${DMC_FREQ}/governor 2>/dev/null || echo dmc_$1 >${DMC_FREQ}/governor 2>/dev/null
}

set_gpu_gov() {
  echo $1 >${GPU_FREQ}/governor 2>/dev/null || echo dmc_$1 >${GPU_FREQ}/governor 2>/dev/null
}

performance() {
  set_cpu_gov performance
  set_gpu_gov performance
  set_dmc_gov performance
}

ondemand() {
  set_cpu_gov ondemand
  set_gpu_gov ondemand
  set_dmc_gov ondemand
}

interactive() {
  set_cpu_gov interactive
  set_gpu_gov ondemand
  set_dmc_gov ondemand
}

powersave() {
  set_cpu_gov powersave
  set_gpu_gov powersave
  set_dmc_gov powersave
}
