# packages/jelos/profile.d/99-freqfunctions
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
# Copyright (C) 2021-present Fewtarius (https://github.com/fewtarius)

. /etc/os-release

set_gov() {
  for POLICY in ${CPU_FREQ[@]}
  do
    if [ "${POLICY}" = "${NO_MAX_PERF}" ] && [ "$1" = "performance" ]
    then
      continue
    else
      echo $1 >${POLICY} 2>/dev/null
    fi
  done

  for HANDLE in ${GPU_FREQ} ${DMC_FREQ}
  do
    echo $1 >${HANDLE} 2>/dev/null || echo dmc_$1 >${HANDLE} 2>/dev/null
  done
}

maxperf() {
  set_gov performance
}

normperf() {
  set_gov ondemand
}

powersave() {
  set_gov powersave
}

