# packages/network/hostapd/system.d/hostapd.service
[Unit]
Description=Hostapd service
After=network.target

[Service]
Environment=XDG_RUNTIME_DIR=/var/run/0-runtime-dir
Environment=HOME=/storage
ExecStart=/usr/bin/hostapd -S -i wlan0 /storage/.cache/.hostapd.conf
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
