# packages/hardware/quirks/profile.d/999-export
#!/bin/bash
# SPDX-License-Identifier: Apache-2.0
# Copyright (C) 2021-present Fewtarius (https://github.com/fewtarius)

### If you add a variable that should persist across processes,
### remember to export it here.

export  SLOW_CORES \
	FAST_CORES \
	CPU_FREQ \
	DMC_FREQ \
	GPU_FREQ \
	DEVICE_AUDIO_MIXER \
	DEVICE_BASE_TDP \
	DEVICE_BRIGHTNESS \
	DEVICE_BATTERY_LED_STATUS \
	DEVICE_FAKE_JACKSENSE \
	DEVICE_FUNC_KEYA_MODIFIER \
	DEVICE_FUNC_KEYB_MODIFIER \
	DEVICE_HAS_FAN \
	DEVICE_HAS_HDMI \
	DEVICE_HDMI_GPIO \
	DEVICE_HEADPHONE_DEV \
	DEVICE_INTERNAL_WIFI \
	DEVICE_JACK \
	DEVICE_KEY_VOLUMEDOWN \
	DEVICE_KEY_VOLUMEUP \
	DEVICE_LED_CONTROL \
        DEVICE_PIPEWIRE_PROFILE \
	DEVICE_PWR_LED_CONTROL \
	DEVICE_PWR_LED_GPIO \
	DEVICE_PLAYBACK_PATH_HP \
	DEVICE_PLAYBACK_PATH_SPK \
	DEVICE_POWER_LED \
	DEVICE_PWM_FAN \
	DEVICE_PWM_MOTOR \
	DEVICE_ROTATED_SDL \
	DEVICE_SW_HP_SWITCH \
	DEVICE_TEMP_SENSOR \
	DEVICE_VOLUME \
	DEVICE_VOLUMECTL \
	DEVICE_WIFI \
	DEVICE_WIFI_MODULE \
	HDMI_DETECT_PATTERN \
	HDMI_STATE \
	UI_SERVICE
