# packages/sysutils/system-utils/udev.d/99-hdmi.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/hdmi_sense"
