# packages/games/tools/rinputer2/system.d/Rinputer2.service
[Unit]
Description=Ragnarok Input Daemon

[Service]
Type=simple
ExecStart=/usr/bin/rinputer2

[Install]
WantedBy=multi-user.target
